# frontend.conf
# yavdr-frontend - Starts the yavdr-frontend
#

description     "frontend"
author          "Lars Hanisch <dvb@flensrocker.de>"

task

start on stopped openbox-tools or frontend-restart
stop on stopping openbox

env HOME=/var/lib/vdr
export HOME

env DISPLAY=:1
export DISPLAY

pre-start script
# dont start if some other application is running
if [ -e /tmp/.standalone ]; then
   /usr/bin/vdr-dbus-send /Remote remote.Disable ||:
   exit 1
fi
end script

script
  /usr/bin/logger -t frontend starting frontend
  if [ -r /etc/default/locale ]; then
        . /etc/default/locale
        export LANG LANGUAGE LC_CTYPE=$LANG
  fi

  exec /usr/bin/yavdr-frontend -c /etc/default/yavdr-frontend
end script
