# unset-session.conf

start on desktop-end or stopping dbus

task

script
  /usr/bin/vdr-dbus-send /Plugins/dbus2vdr plugin.SVDRPCommand string:"SetDBusSessionBusAddress" string:"" > /dev/null
  /usr/bin/vdr-dbus-send /Plugins/dbus2vdr plugin.SVDRPCommand string:"SetUpstartSession" string:"" > /dev/null
end script

