# var/lib/vdr/.config/upstart/on-vdr.conf
start on started dbus
stop on stopping dbus

script
  /usr/bin/on_vdr --one-shot -c "/sbin/initctl --user emit vdr-start"
end script

