# var/lib/vdr/.config/upstart/udisks-automounter.conf

description     "start udisks-automounter to automount and notify vdr"

start on startup
stop on desktop-end

env DISPLAY=:1

export DISPLAY

script
exec ck-launch-session udisks-glue -c /etc/udisks-glue/config --foreground
end script

post-stop script
#killall udisks-glue || /bin/true
end script

